# root/etc/systemd/system/pve-guests.service.d/manual-start.conf
[Unit]
RefuseManualStart=false
RefuseManualStop=false

